# deploy/tgbot.service.example
# Пример unit для systemd (Linux VPS).
# Проще включить автозапуск скриптом из корня репозитория на VPS:
#   sudo bash deploy/install-systemd.sh /home/ubuntu/bots/TGBOTAImbp77 ubuntu
# Или скопируйте этот файл в /etc/systemd/system/tgbot-imbp77.service,
# поправьте User / WorkingDirectory / ExecStart, затем:
#   sudo systemctl daemon-reload && sudo systemctl enable --now tgbot-imbp77

[Unit]
Description=TGBOTAImbp77 Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/bots/TGBOTAImbp77
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/bots/TGBOTAImbp77/.venv/bin/python /home/ubuntu/bots/TGBOTAImbp77/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
