# deploy/tgbot.service.example
# Пример unit для systemd (Linux VPS).
# Скопируйте в /etc/systemd/system/tgbot-imbp77.service, поправьте User и WorkingDirectory,
# затем: sudo systemctl daemon-reload && sudo systemctl enable --now tgbot-imbp77

[Unit]
Description=TGBOTAImbp77 Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/TGBOTAImbp77
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/TGBOTAImbp77/.venv/bin/python /opt/TGBOTAImbp77/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
