# deploy/tgbot.service.example
# Пример unit для systemd (Linux VPS).
# Поправьте User, WorkingDirectory и ExecStart под каталог, куда сделали git clone
# (например /home/ВАШ_ЮЗЕР/bots/TGBOTAImbp77 или /opt/tgbot-imbp77).
# Затем: sudo systemctl daemon-reload && sudo systemctl enable --now tgbot-imbp77

[Unit]
Description=TGBOTAImbp77 Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/bots/TGBOTAImbp77
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/bots/TGBOTAImbp77/.venv/bin/python /home/ubuntu/bots/TGBOTAImbp77/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
